# examples/etc/supervisor.d/overpush.ini
[program:overpush]
command=/usr/local/bin/overpush
user=www
directory=/var/www/htdocs
stdout_syslog=true
stderr_syslog=true
